# scripts/com.xiaoyibanyun.watch-openapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xiaoyibanyun.watch-openapi</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/sunsh80/Downloads/易工到项目/小蚁搬运/scripts/watch-openapi.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sunsh80/Downloads/易工到项目/小蚁搬运</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/sunsh80/Downloads/易工到项目/小蚁搬运/scripts/watch-openapi-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sunsh80/Downloads/易工到项目/小蚁搬运/scripts/watch-openapi-stderr.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
